# templates/etc/systemd/debian_redis-server.service.erb
[Unit]
Description=Redis Datastore Server
After=network.target

[Service]
Type=forking
PIDFile=<%= @redis_pid_dir %>/redis_<%= @redis_name %>.pid
ExecStartPre=/bin/mkdir -p <%= @redis_pid_dir %>
ExecStartPre=/bin/chown redis:redis <%= @redis_pid_dir %>
ExecStartPre=/bin/mkdir -p <%= @redis_log_dir %>
ExecStartPre=/bin/chown redis:redis <%= @redis_log_dir %>

ExecStart=/sbin/start-stop-daemon --start --chuid redis:redis --pidfile <%= @redis_pid_dir %>/redis_<%= @redis_name %>.pid --umask 007 --exec <%= @redis_install_dir %>/redis-server -- /etc/redis_<%= @redis_name %>.conf
ExecReload=/bin/kill -USR2 $MAINPID

[Install]
WantedBy=multi-user.target
