# templates/systemd/sentinel.service.erb
[Unit]
Description=Redis Sentinel <%= @sentinel_name %>
After=network.target

[Service]
<% if @sentinel_user and @sentinel_user != 'root' -%>
PermissionsStartOnly=true
ExecStartPre=/bin/mkdir -p <%= @sentinel_run_dir %>
ExecStartPre=/bin/chown <%= @sentinel_user %>:<%= @sentinel_group or 'root' %> <%= @sentinel_run_dir %>
<% else -%>
ExecStartPre=/bin/mkdir -p <%= @sentinel_run_dir %>
<% end -%>
ExecStartPre=/bin/cp -u <%= @conf_file %> <%= @sentinel_run_dir %>/<%= @conf_file_name %>
ExecStartPre=/bin/chown <%= @sentinel_user or 'root' %>:<%= @sentinel_group or 'root' %> <%= @sentinel_run_dir %>/<%= @conf_file_name %>
ExecStart=/usr/bin/redis-sentinel  <%= @sentinel_run_dir %>/<%= @conf_file_name %> --daemonize no
ExecStop=/usr/bin/redis-cli -p <%= @sentinel_port %> shutdown
User=<%= @sentinel_user or 'root' %>
Group=<%= @sentinel_group or 'root' %>

[Install]
WantedBy=multi-user.target

